# portal/config/cron.d/portal
* *  * * *  root /opt/Portal-v4/portal/watchdog.sh
* *  * * *  root /opt/Portal-v4/portal/send.py > /dev/null
